# hpsu.medium.service
[Unit]
Description=HPSU MEDIUM Daemon

[Service]
Type=simple
ExecStart=/home/pi/hpsu.medium.sh
PIDFile=/tmp/hpsu.medium.pid
RuntimeMaxSec=8
