# hpsu.slow.service
[Unit]
Description=HPSU SLOW Daemon

[Service]
Type=simple
ExecStart=/home/pi/hpsu.slow.sh
PIDFile=/tmp/hpsu.slow.pid
RuntimeMaxSec=8
