# scripts/capture-publish/hpsu.listener.service
[Unit]
Description=HPSU LISTENER Daemon

[Service]
Type=simple
ExecStart=/home/pi/hpsu.listener.sh
PIDFile=/tmp/hpsu.listener.pid
